# systemd/familyhub.service
[Unit]
Description=FamilyHub Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# To enable sdnotify/watchdog later, switch to:
# Type=notify
# WatchdogSec=30s
User=familyhub
WorkingDirectory=/home/pi/familyhub
EnvironmentFile=/etc/default/familyhub
ExecStart=/home/pi/familyhub/.venv/bin/python /home/pi/familyhub/main.py
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
